# systemd/mobile-terminal.service
[Unit]
Description=Mobile Terminal
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=@WORKDIR@
EnvironmentFile=-@ENV_FILE@
ExecStart=@PYTHON@ @WORKDIR@/server.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
